# django_examples/store/supervisor.conf
[program:store]
command=/home/ubuntu/store/gunicorn_start
directory=/home/ubuntu/store
user=ubuntu
stdout_logfile=/home/ubuntu/logs/gunicorn_supervisor.log
redirect_stderr=true
